# systemd/vocera-media-qoe-wlc-session-ingest.service
[Unit]
Description=Ingest exported WLC capture-session EPCs SCP-pushed to the collector
Documentation=file:/home/appsadmin/grafana-mimir-observability/docs/wireless/vocera-media-dnac-icap-runbook.md
After=network-online.target vocera-rf-validation-study-web.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/appsadmin/grafana-mimir-observability
Environment=STUDY_WEB_INGEST_HOST=127.0.0.1
Environment=STUDY_WEB_INGEST_PORT=8097
Environment=STUDY_WEB_INGEST_TIMEOUT=600
# Required ordering so a valid large-EPC parse is never killed mid-flight:
#   TimeoutStartSec (660) > curl max-time / STUDY_WEB_INGEST_TIMEOUT (600)
#     > Study Web parser timeout STUDY_WEB_MEDIA_QOE_PARSE_TIMEOUT_SECONDS (480).
TimeoutStartSec=660
EnvironmentFile=-/etc/default/vocera-media-qoe-wlc-session-ingest
ExecStart=/usr/bin/bash scripts/run_vocera_wlc_session_ingest.sh
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
