# systemd/wireless-path-probe.timer
[Unit]
Description=Run wireless path RTT probes every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=wireless-path-probe.service

[Install]
WantedBy=timers.target
