# systemd-overrides/grafana-server.service.d/secrets.conf
# systemd drop-in for grafana-server. Loads per-service Postgres passwords
# written by scripts/install_secrets.sh so the Grafana datasource yamls'
# ${TOPOLOGY_POSTGRES_PASSWORD} etc. expansions resolve at provisioning time.
#
# The `-` prefix on each EnvironmentFile= makes the file optional, so Grafana
# still starts (with the matching datasource broken) if a particular secrets
# file is missing. This keeps a partial install from blocking grafana-server.

[Service]
EnvironmentFile=-/etc/grafana-mimir-observability/secrets/topology-postgres.env
EnvironmentFile=-/etc/grafana-mimir-observability/secrets/vocera-media-qoe-postgres.env
EnvironmentFile=-/etc/grafana-mimir-observability/secrets/vocera-rf-validation-postgres.env
